# ee/appliance/ubuntu-iso/overlay/etc/systemd/system/alga-appliance.service
[Unit]
Description=Alga Appliance Host Setup and Status Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=ALGA_APPLIANCE_PORT=8080
Environment=ALGA_APPLIANCE_STATE_FILE=/var/lib/alga-appliance/install-state.json
Environment=ALGA_APPLIANCE_TOKEN_FILE=/var/lib/alga-appliance/setup-token
ExecStartPre=/usr/bin/env node /opt/alga-appliance/host-service/init-token.mjs
ExecStart=/usr/bin/env node /opt/alga-appliance/host-service/server.mjs
Restart=always
RestartSec=3
WorkingDirectory=/opt/alga-appliance
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
